# init_router.service
[Unit]
After=network.service

[Service]
ExecStart=/script/init_router.sh

[Install]
WantedBy=default.target
